# batch-processor/crontab
0 20 * * * docker exec batch-processor python /app/batch_processor.py

